# train-spam-scanner.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# merge user trained messages into admin training folders
35 * * * *  root  /usr/local/sbin/train-spam-scanner -F -M none

# train the spam scanner with sorted mail
15 3 * * *  root  /usr/local/sbin/train-spam-scanner -f
